# includes.container/etc/systemd/system/waydroid-init.service
[Unit]
Description=Initialize the Waydroid container

[Service]
User=root
WorkingDirectory=/root
ExecStart=waydroid init

[Install]
WantedBy=multi-user.target
